# examples/systemd/pubsub.service
Description=Crawl PubSub

[Service]
ExecStart=/app/pubsub/examples/run.sh

[Install]
WantedBy=multi-user.target
